# RPi_Motor_Control/rpi-motor-control.service
[Unit]
Description=RPi Motor Control System
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/omar-drill-machine-project-rpi-python/RPi_Motor_Control
Environment="PATH=/home/ubuntu/omar-drill-machine-project-rpi-python/RPi_Motor_Control/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/omar-drill-machine-project-rpi-python/RPi_Motor_Control/venv/bin/python3 /home/ubuntu/omar-drill-machine-project-rpi-python/RPi_Motor_Control/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
