# RPi_Motor_Control/scripts/motor_control.service
[Unit]
Description=RPi Motor Control System
After=network.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/RPi_Motor_Control
Environment="PATH=/home/pi/RPi_Motor_Control/venv/bin"
ExecStart=/home/pi/RPi_Motor_Control/venv/bin/python /home/pi/RPi_Motor_Control/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
